# filesystem/etc/systemd/system/movie-daemon.service
[Unit]
Description=MovieScoutApi

[Service]

WorkingDirectory={{ROOT_DIR}}
ExecStart=/usr/bin/yarn run dev

[Install]
WantedBy=multi-user.target

